# configs/pipeline.yaml
# configs/pipeline.yaml
# ─────────────────────────────────────────────────────────────────
# Master configuration — every tunable parameter lives here.
# ─────────────────────────────────────────────────────────────────

paths:
  raw_dir:       "data/raw"
  processed_dir: "data/processed"
  graphs_dir:    "data/graphs"
  splits_dir:    "data/splits"

# ── DATA SOURCES ─────────────────────────────────────────────────
sources:
  otrf:
    base_url: "https://raw.githubusercontent.com/OTRF/Security-Datasets/master/datasets/atomic/windows"
    datasets:
      # PowerShell Empire C2 lateral movement
      - path: "lateral_movement/host/empire_psexec_pth_dcerpc_svcctl_RemoteCreateService.zip"
        ttp:  ["T1021.002", "T1059.001"]
      # Credential access via Mimikatz through PS
      - path: "credential_access/host/empire_mimikatz_extract_keys.zip"
        ttp:  ["T1003.001", "T1059.001"]
      # Defense evasion — encoded commands
      - path: "defense_evasion/host/empire_psinject.zip"
        ttp:  ["T1059.001", "T1027"]
      # Discovery via PowerShell
      - path: "discovery/host/empire_find_local_admin_access.zip"
        ttp:  ["T1069.001", "T1059.001"]
      # Persistence — scheduled task via PS
      - path: "persistence/host/empire_schtasks.zip"
        ttp:  ["T1053.005", "T1059.001"]

  splunk_attack_data:
    base_url: "https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques"
    datasets:
      # T1059.001 — PowerShell execution
      - path: "T1059.001/atomic_red_team/windows-sysmon.log"
        ttp:  ["T1059.001"]
      # T1027 — obfuscated files
      - path: "T1027/atomic_red_team/windows-sysmon.log"
        ttp:  ["T1027", "T1059.001"]
      # T1003.001 — LSASS memory dump via PS
      - path: "T1003.001/atomic_red_team/windows-sysmon.log"
        ttp:  ["T1003.001", "T1059.001"]

# ── SYSMON EVENT IDs TO COLLECT ──────────────────────────────────
sysmon:
  event_ids:
    - 1    # Process creation (primary — has CommandLine, ParentImage)
    - 3    # Network connection
    - 7    # Image loaded
    - 10   # Process accessed (lsass injection signals)
    - 11   # File created
    - 17   # Pipe created
    - 18   # Pipe connected
  powershell_images:
    - "powershell.exe"
    - "pwsh.exe"
  scope_hops: 3   # Follow process tree this many hops from PS root

# ── TABULAR FEATURES (RF Tier 1) ─────────────────────────────────
features:
  tabular:
    # Structural
    - has_encoded_arg          # -enc or -EncodedCommand present
    - has_noprofile            # -nop / -NoProfile
    - has_hidden_window        # -w hidden / -WindowStyle Hidden
    - has_bypass               # -ExecutionPolicy Bypass
    - has_iex                  # Invoke-Expression / iex
    - has_downloadstring       # Net.WebClient / DownloadString / IWR
    - has_reflection           # [Reflection.Assembly] / Add-Type
    - has_wmi                  # WMI / gwmi / Invoke-WmiMethod
    - has_scheduled_task       # schtasks / Register-ScheduledTask
    - has_registry_write       # Set-ItemProperty / reg add / HKLM
    - has_credential_access    # mimikatz / lsass / sekurlsa
    - has_lateral_movement     # psexec / Enter-PSSession / Invoke-Command + ComputerName
    - has_network_call         # any Event ID 3 child in session
    - has_pipe_activity        # Event ID 17/18 in session
    - cmd_length               # raw character count of CommandLine
    - entropy_score            # Shannon entropy of CommandLine
    - parent_is_office         # ParentImage is winword/excel/outlook/powerpnt
    - parent_is_browser        # ParentImage is chrome/firefox/msedge/iexplore
    - parent_is_script_host    # ParentImage is wscript/cscript/mshta
    - parent_is_service        # ParentImage is services.exe/svchost.exe
    - parent_is_powershell     # PS spawned from another PS (nested)
    - user_is_system           # IntegrityLevel = System
    - user_is_high             # IntegrityLevel = High
    - child_count              # number of child processes spawned
    - unique_child_images      # distinct child process images
    - session_event_count      # total events in the session
    - session_duration_secs    # time from first to last event in session
    - has_outbound_443         # network connection on port 443 (C2 beacon)
    - has_outbound_80          # network connection on port 80
    - has_nonstandard_port     # outbound on ports not in [80, 443, 53, 8080]
    - file_write_count         # Event ID 11 count in session
    - has_temp_write           # file write to %TEMP% or %APPDATA%
    - has_system32_spawn       # child process image in System32
    - base64_token_count       # count of Base64-like tokens (len>20, alphanum+/+=)
    - obfuscation_score        # ratio of non-printable / special chars in CLI
    - token_count              # whitespace-split token count of CLI
    - unique_cmdlet_count      # count of distinct PS cmdlets detected
    - has_double_extension     # file args with double extension (doc.exe etc)
    - lolbin_child_count       # children that are known LOLBins
    - script_block_length      # from Event 4104 if available, else 0
    - depth_in_tree            # depth of this node from session root

# ── GRAPH SCHEMA (GAT Tier 2B) ───────────────────────────────────
graph:
  node_feature_dim: 18
  node_features:
    - is_powershell            # binary: image is powershell.exe or pwsh.exe
    - is_lolbin                # binary: image is in LOLBAS list
    - integrity_level          # 0=low, 1=medium, 2=high, 3=system
    - has_network_connection   # binary: Event ID 3 exists for this PID
    - has_pipe_activity        # binary: Event ID 17/18 exists for this PID
    - has_file_write           # binary: Event ID 11 exists for this PID
    - cmd_entropy              # float: Shannon entropy of CommandLine
    - cmd_length_norm          # float: log(len+1) normalised to [0,1]
    - has_encoded_cmd          # binary
    - has_download_pattern     # binary
    - has_credential_pattern   # binary
    - parent_depth             # int: hops from session root
    - child_count              # int: number of direct children
    - user_privilege           # 0=low, 1=medium, 2=high, 3=system
    - time_delta_secs          # float: seconds since previous event in session
    - image_hash               # int: hash(image_name) mod 1000 (process identity)
    - has_wmi_activity         # binary
    - has_lateral_pattern      # binary

  edge_types:
    parent_spawn: 0            # directed: parent → child (Event ID 1 ParentPID)
    pipe_connect: 1            # undirected: processes sharing a pipe
    file_write:   2            # directed: writer → file node (simplified)
    net_connect:  3            # directed: process → remote endpoint

# ── LABELLING ────────────────────────────────────────────────────
labelling:
  # A session is malicious if ANY event in it is attack-sourced
  strategy: "session_level"
  # TTP tags are inherited from dataset metadata (not heuristic)
  ttp_source: "dataset_metadata"
  # Minimum events to include a session (filter noise)
  min_session_events: 3
  # Max session gap — if >N seconds between PS events, start new session
  session_gap_secs: 300

# ── TRAIN / VAL / TEST SPLIT ─────────────────────────────────────
split:
  train: 0.70
  val:   0.10
  test:  0.20
  stratify_by: "label"        # ensure class balance across splits
  random_seed: 42
  # Hold out entire attack families for OOD test
  ood_families:
    - "empire_mimikatz_extract_keys"   # credential access family
